# archiveofsound@.service
[Unit]
Description=Archive of Sound

[Service]
Type=exec
ExecStart=archiveofsound %i
